# bridge/systemd/qabridge.service
[Unit]
Description=qabridge appium service

[Service]
User=root
ExecStart=/usr/bin/qabridge

[Install]
WantedBy=multi-user.target
